# etc/init.d/bcsmd
#!/sbin/openrc-run

description="BroadCast Short Messages Daemon"
command="${MSGD_HOME}/env/bcsmd.sh"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

