# etc/init.d/msgd
#!/sbin/openrc-run

description="Message daemon"
command="${MSGD_HOME}/env/msgd.sh"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"

