# etc/init.d/telegram-bot-bash
#!/sbin/openrc-run

TGBB_DSC="telegram bot bash"

run_bot_cmd()
{
	TGBB_CMD="cd ${TGBB_HOME} && ./bashbot.sh $@"
	su - ${TGBB_USER} -c "${TGBB_SHELL} -lc \"${TGBB_CMD}\"" >/dev/null
}

depend()
{
	need net
}

start()
{
	run_bot_cmd start

	ebegin "Start $TGBB_DSC"

	run_bot_cmd status && eend 0 || eend 1
}

stop()
{
	run_bot_cmd stop

	ebegin "Stop $TGBB_DSC"

	run_bot_cmd status && eend 1 || eend 0
}

